# profile.d/awshelper.sh
function awssso() {
    aws sso login --profile lab
    aws sso login --profile prod
}
